# initramfs/etc/inittab
::sysinit:/bin/busybox --install
::sysinit:/bin/busybox mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/busybox mount -t proc proc /proc
/dev/console::sysinit:-/bin/sh
